# tools/systemd/psyche-client@.service
[Unit]
Description=Psyche Training Client (%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/src/psyche
EnvironmentFile=%h/.config/psyche/client-%i.env
ExecStart=/bin/bash -c 'exec cargo run --release --bin psyche-solana-client -- train \
    --wallet-private-key-path ${WALLET_PRIVATE_KEY_PATH} \
    --rpc ${RPC} \
    --ws-rpc ${WS_RPC} \
    --run-id %i \
    --data-parallelism ${DP:-1} \
    --tensor-parallelism ${TP:-1} \
    --micro-batch-size ${BATCH_SIZE:-1} \
    --authorizer ${AUTHORIZER:-11111111111111111111111111111111} \
    --logs console \
    ${EXTRA_ARGS:-}'
Environment=RUST_LOG=info,psyche=debug
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
